# templates/systemd/wt-web.service
[Unit]
Description=wt-tools Web Dashboard
Documentation=https://github.com/tatargabor/wt-tools
After=default.target

[Service]
ExecStart=%h/.local/bin/wt-orch-core serve --port 7400
Restart=always
RestartSec=5
Environment=PYTHONPATH=%h/code2/wt-tools/lib

[Install]
WantedBy=default.target
